# home/private_dot_config/exact_ansible/exact_roles/exact_linux/exact_applications/exact_home-assistant/exact_files/docker.home-assistant.service
[Unit]
Description=HomeAssistant Docker Container through Podman

[Service]
Restart=always
ExecStart=/usr/bin/podman start -a home-assistant
ExecStop=/usr/bin/podman stop -t 2 home-assistant

[Install]
WantedBy=multi-user.target
